# scripts/deploy/brainstem.service
[Unit]
Description=Brainstem CMS gRPC Server (medulla)
After=network.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/data/brainstem
Environment=MEDULLA_STANDALONE=1
Environment=MEDULLA_PORT=13145
Environment=MEDULLA_LOG=INFO
Environment=PATH=/home/sunrise/data/dart-sdk/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/sunrise/data/dart-sdk/bin/dart run han_dog/bin/server.dart
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
