# scripts/deploy/gnss.service
[Unit]
Description=LingTu GNSS Driver (WTRTK-980 / UM980)
Documentation=https://github.com/ironoa/um982_ros2_driver
After=network.target
Wants=network.target

# Camera 和 livox-driver 不依赖 GNSS，但我们希望 GNSS 先启动好
# 让 lingtu.service 启动时能立即订阅 /gps/fix
Before=lingtu.service

[Service]
Type=simple
User=sunrise
Group=dialout

# Device permissions: sunrise 必须在 dialout 组才能读 /dev/ttyUSB*
Environment=ROS_DOMAIN_ID=0
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp

# wtrtk980_ros2_reader — C++ ROS2 node for WTRTK-980
# Source: src/drivers/gnss/wtrtk980_ros2_reader/
# Build:  ~/data/inovxio/lingtu/install/ (colcon)
ExecStartPre=/bin/bash -c 'test -e /dev/wtrtk980'
ExecStart=/bin/bash -c " \
    source /opt/ros/humble/setup.bash && \
    source /home/sunrise/data/inovxio/lingtu/install/setup.bash && \
    ros2 run wtrtk980_ros2_reader wtrtk980_ros2_node \
      --ros-args \
      -p device:=auto \
      -p fix_topic:=/gps/fix \
      -p raw_topic:=/gps/nmea \
      -p frame_id:=wtrtk980 \
      -p publish_raw:=true"

Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=control-group
KillSignal=SIGTERM

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gnss

[Install]
WantedBy=multi-user.target
