# scripts/deploy/lingtu-manager.service
[Unit]
Description=LingTu Manager — Robot Control Panel
After=network.target lidar.service camera.service
Wants=lidar.service camera.service

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/data/inovxio/lingtu
Environment=PYTHONUNBUFFERED=1
Environment=NAV_MAP_DIR=/home/sunrise/data/inovxio/data/maps
ExecStart=/usr/bin/python3 scripts/manager/manager.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
