# scripts/deploy/localizer.service
[Unit]
Description=LingTu ICP Localizer
After=slam.service

[Service]
Type=simple
User=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash -c 'source /opt/nav/scripts/services/env.sh && exec ros2 launch /opt/nav/launch/profiles/localizer_icp.launch.py'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
