# scripts/deploy/slam.service
[Unit]
Description=LingTu SLAM (Fast-LIO2)
After=lidar.service

[Service]
Type=simple
User=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash -c 'source /opt/nav/scripts/services/env.sh && exec ros2 run fastlio2 lio_node --ros-args -p config_path:=${NAV_DIR}/install/fastlio2/share/fastlio2/config/lio.yaml -r /cloud_registered:=/nav/registered_cloud -r /Odometry:=/nav/odometry -r /cloud_map:=/nav/map_cloud -r /imu/data:=/nav/imu -r /lidar/scan:=/nav/lidar_scan'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
