# scripts/deploy/slam_pgo.service
[Unit]
Description=LingTu PGO (included in SLAM fastlio2 profile — this service runs PGO standalone)
After=slam.service

[Service]
Type=simple
User=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash -c 'source /opt/nav/scripts/services/env.sh && exec ros2 run pgo pgo_node --ros-args -p config_path:=${NAV_DIR}/install/pgo/share/pgo/config/pgo.yaml -r /cloud_registered:=/nav/registered_cloud -r /Odometry:=/nav/odometry'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
