# systemd/nav-monitor.service
[Unit]
Description=Navigation - Three-Ring Monitor + Health HTTP Server
After=network.target

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-monitor.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
KillMode=control-group
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

# 监控服务内存限制宽松
MemoryMax=1G

Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-monitor

[Install]
WantedBy=multi-user.target
