# docs/04-deployment/services/lingtu.service
[Unit]
Description=LingTu Navigation System (release-mode)
Documentation=https://github.com/Kitjesen/MapPilot
# Depend on the SLAM + camera systemd services that publish /nav/odometry
# and /camera/* topics. brainstem is intentionally NOT here — han_dog.service
# is the actual brainstem controller; robot-brainstem.service is a legacy
# duplicate that has been masked.
After=robot-fastlio2.service robot-localizer.service robot-camera.service network.target
Wants=robot-fastlio2.service robot-localizer.service robot-camera.service

[Service]
Type=simple
User=sunrise

# Release-mode: read code + binaries from the immutable release dir.
# /opt/lingtu/current is a symlink that scripts/deploy/cut_release.sh swaps
# atomically when a new version is cut.
WorkingDirectory=/opt/lingtu/current

Environment=PYTHONUNBUFFERED=1

# DDS workaround — cyclonedds Topic() segfaults on S100P when subscribing
# to a fastrtps-published topic (QoS mismatch at C layer). Forcing rclpy
# fallback avoids the crash. Remove once the root cause is fixed.
Environment=LINGTU_DISABLE_DDS=1

# ExecStart: bash wrapper sources ROS env + user API keys, then exec lingtu.py.
#   1. /opt/ros/humble/setup.bash       — ROS2 core (AMENT_PREFIX_PATH, LD_LIBRARY_PATH)
#   2. install/setup.bash               — colcon-built workspace (livox/fastlio2/etc.)
#   3. eval API keys from ~/.bashrc     — DASHSCOPE / OPENAI / MOONSHOT / ANTHROPIC
#                                         (LLMModule preflight requires at least one)
#   4. exec lingtu.py nav --no-repl     — production navigation profile, no REPL
ExecStart=/bin/bash -c ' \
  set -e; \
  source /opt/ros/humble/setup.bash; \
  source /opt/lingtu/current/install/setup.bash; \
  eval "$(grep -E "^export (DASHSCOPE_API_KEY|OPENAI_API_KEY|MOONSHOT_API_KEY|ANTHROPIC_API_KEY|LINGTU_API_KEY)=" /home/sunrise/.bashrc || true)"; \
  exec /usr/bin/python3 /opt/lingtu/current/lingtu.py nav --no-repl'

Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=15s

[Install]
WantedBy=multi-user.target
