# docs/04-deployment/services/robot-camera.service
[Unit]
Description=Robot Camera Driver (Orbbec Gemini 330)
Documentation=https://github.com/Kitjesen/lingtu
After=network.target

[Service]
Type=simple
User=sunrise
ExecStart=/bin/bash -c '\
  source /opt/ros/humble/setup.bash && \
  exec ros2 launch orbbec_camera gemini_330_series.launch.py'
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s

[Install]
WantedBy=multi-user.target
