# docs/04-deployment/services/robot-lidar.service
[Unit]
Description=Robot LiDAR Driver (Livox MID360)
Documentation=https://github.com/Kitjesen/lingtu
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/opt/lingtu/current
ExecStartPre=/bin/bash -c 'source /opt/lingtu/config/ros2-env.sh && ros2 daemon stop 2>/dev/null; true'
ExecStart=/bin/bash -c '\
  source /opt/lingtu/config/ros2-env.sh && \
  exec /opt/lingtu/current/install/livox_ros_driver2/lib/livox_ros_driver2/livox_ros_driver2_node \
    --ros-args \
    -p xfer_format:=1 \
    -p multi_topic:=0 \
    -p data_src:=0 \
    -p publish_freq:=10.0 \
    -p output_data_type:=0 \
    -p user_config_path:=/opt/lingtu/current/install/livox_ros_driver2/share/livox_ros_driver2/config/MID360_config.json \
    -r livox/lidar:=/nav/lidar_scan \
    -r livox/imu:=/nav/imu'
ExecStartPost=/bin/bash -c '\
  source /opt/lingtu/config/ros2-env.sh && \
  for _ in $(seq 1 15); do \
    lidar=$$(ros2 topic info /nav/lidar_scan 2>/dev/null | sed -n "s/^Publisher count: //p" | tail -1); \
    imu=$$(ros2 topic info /nav/imu 2>/dev/null | sed -n "s/^Publisher count: //p" | tail -1); \
    if [ "$${lidar:-0}" -ge 1 ] && [ "$${imu:-0}" -ge 1 ]; then exit 0; fi; \
    sleep 1; \
  done; \
  echo "robot-lidar readiness failed: /nav/lidar_scan=$${lidar:-0} /nav/imu=$${imu:-0}" >&2; \
  exit 1'
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=10s
KillMode=control-group
SendSIGKILL=yes

[Install]
WantedBy=multi-user.target
