# scripts/deploy/s100p/lidar.service
[Unit]
Description=S100P Livox MID-360 LiDAR Driver
After=network.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/data/SLAM/navigation
ExecStart=/bin/bash -c '\
  source /opt/ros/humble/setup.bash && \
  source ~/data/SLAM/navigation/install/setup.bash && \
  exec ros2 run livox_ros_driver2 livox_ros_driver2_node \
    --ros-args \
    -p xfer_format:=1 \
    -p multi_topic:=0 \
    -p data_src:=0 \
    -p publish_freq:=10.0 \
    -p output_data_type:=0 \
    -p user_config_path:=/home/sunrise/data/SLAM/navigation/install/livox_ros_driver2/share/livox_ros_driver2/config/MID360_config.json \
    -r livox/lidar:=/nav/lidar_scan \
    -r livox/imu:=/nav/imu'
Restart=on-failure
RestartSec=3
MemoryMax=1G
KillMode=control-group
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lidar
Environment=HOME=/home/sunrise
Environment=ROS_DOMAIN_ID=0

[Install]
WantedBy=multi-user.target
