# scripts/deploy/s100p/slam_pgo.service
[Unit]
Description=S100P PGO (Pose Graph Optimization)
After=slam.service
Wants=slam.service

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/data/inovxio/lingtu
ExecStart=/bin/bash -c '\
  source /opt/ros/humble/setup.bash && \
  source ~/data/inovxio/lingtu/install/setup.bash && \
  exec ros2 run pgo pgo_node \
    --ros-args \
    -p config_path:=/home/sunrise/data/inovxio/lingtu/install/pgo/share/pgo/config/pgo.yaml \
    -r /cloud_registered:=/nav/registered_cloud \
    -r /Odometry:=/nav/odometry'
Restart=on-failure
RestartSec=3
MemoryMax=2G
KillMode=control-group
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slam_pgo
Environment=HOME=/home/sunrise
Environment=ROS_DOMAIN_ID=0

[Install]
WantedBy=multi-user.target
