# scripts/deploy/s100p/super_lio.service
[Unit]
Description=S100P Super-LIO experimental LIO backend
After=network.target lidar.service robot-lidar.service
Wants=lidar.service robot-lidar.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
TimeoutStopSec=15
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/data/inovxio/lingtu
Environment=HOME=/home/sunrise
Environment=ROS_DOMAIN_ID=0
Environment=SUPER_LIO_WS=/home/sunrise/data/inovxio/super-lio
Environment=SUPER_LIO_LIDAR_TOPIC=/nav/lidar_scan
Environment=SUPER_LIO_IMU_TOPIC=/nav/imu
EnvironmentFile=-/etc/lingtu/super_lio.env
ExecStart=/bin/bash -c '\
  set -eo pipefail; \
  source /opt/ros/humble/setup.bash; \
  LINGTU_LIVOX_SETUP="/opt/lingtu/current/install/livox_ros_driver2/share/livox_ros_driver2/local_setup.bash"; \
  if [ -f "$${LINGTU_LIVOX_SETUP}" ]; then \
    source "$${LINGTU_LIVOX_SETUP}"; \
  elif [ -f "/opt/lingtu/current/install/setup.bash" ]; then \
    source "/opt/lingtu/current/install/setup.bash" || true; \
  fi; \
  if [ -f "$${SUPER_LIO_WS}/install/local_setup.bash" ]; then \
    source "$${SUPER_LIO_WS}/install/local_setup.bash"; \
  elif [ -f "$${SUPER_LIO_WS}/install/setup.bash" ]; then \
    source "$${SUPER_LIO_WS}/install/setup.bash"; \
  elif [ -f "/home/sunrise/data/inovxio/Super-LIO/install/local_setup.bash" ]; then \
    source "/home/sunrise/data/inovxio/Super-LIO/install/local_setup.bash"; \
  elif [ -f "/home/sunrise/data/inovxio/Super-LIO/install/setup.bash" ]; then \
    source "/home/sunrise/data/inovxio/Super-LIO/install/setup.bash"; \
  else \
    echo "Super-LIO install setup not found. Set SUPER_LIO_WS in /etc/lingtu/super_lio.env." >&2; \
    exit 1; \
  fi; \
  CONFIG="$${SUPER_LIO_CONFIG:-$$(ros2 pkg prefix super_lio)/share/super_lio/config/livox_360.yaml}"; \
  exec ros2 run super_lio super_lio_node \
    --ros-args \
    --params-file "$${CONFIG}" \
    -p "lio.ros.lidar_topic:=$${SUPER_LIO_LIDAR_TOPIC}" \
    -p "lio.ros.imu_topic:=$${SUPER_LIO_IMU_TOPIC}" \
    -r "/lio/odom:=/nav/odometry" \
    -r "/lio/cloud_world:=/nav/map_cloud" \
    -r "/lio/imu/odom:=/nav/imu_odom" \
    -r "/lio/robo/odom:=/nav/robot_odom"'
Restart=on-failure
RestartSec=3
MemoryMax=2G
KillMode=control-group
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=super_lio

[Install]
WantedBy=multi-user.target
